# server/minecraft.service
[Unit]
Description=Vanilla Minecraft server service
After=network.target
StartLimitAction=reboot

[Service]
Type=simple
Restart=on-failure
User=root
ExecStart=/usr/games/minecraft/run.sh
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
